# postfrei.service
[Unit]
Description=Service handler for postfrei

[Service]
Type=simple
ExecStart=/home/myuser/postfrei/app.php

[Install]
WantedBy=multi-user.target
